# templates/default/profile.d/ruby.sh.erb
export PATH=<%= @ruby_bin_path %>:$PATH

